# scripts/assets/gw-warm-quiesce-helper.service
[Unit]
Description=Golem Workers warm-image filesystem quiesce helper
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/gw-warm-quiesce-helper.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
